# Database_Services/billing_and_payment/billing_payment.service
[Unit]
Description=My Python Background Service

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/Dental-Modern-Data-demo/Database_Services/billing_and_payment
ExecStart=/usr/bin/python /home/ubuntu/Dental-Modern-Data-demo/Database_Services/billing_and_payment/bg_worker.py
Restart=always
User=ubuntu

[Install]
WantedBy=multi-user.target