# Database_Services/patient_crm/patient_erm.service
[Unit]
Description=My Python Background Service

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/Dental-Modern-Data-demo/Database_Services/patient_crm
ExecStart=/usr/bin/python /home/ubuntu/Dental-Modern-Data-demo/Database_Services/patient_crm/bg_worker.py
Restart=always
User=ubuntu

[Install]
WantedBy=multi-user.target
